# service/bitlbee/run
#!/bin/bash

set -euo pipefail

exec /usr/sbin/bitlbee -i 127.0.0.1 -p 6667 -P /var/run/bitlbee.pid -F -n
